# www-apps/up-common-proxies/files/up-common-proxies.initd
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="up-common-proxies"
description="UnifiedPush common proxies (bin)"
command="/usr/bin/up-rewrite"
pidfile="/run/${name}.pid"
logdir="/var/log/${name}"
command_background="yes"
command_user="gotify:gotify"
# directory="/var/lib/${name}"
output_log="${logdir}/${name}.log"
error_log="${logdir}/error.log"
command_args="-c /etc/${name}/config.toml"

depend() {
	need net
	use dns
	after logger
}
